# scripts/com.investpilot.noon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investpilot.noon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sjh/Desktop/A股AI投研系统/scripts/noon_push.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/sjh/Desktop/A股AI投研系统/data/logs/launchd_noon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sjh/Desktop/A股AI投研系统/data/logs/launchd_noon.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
